# cmd/flapperd/systemd/flapperd.service
[Unit]
Description=Splitflap Display Control Service

[Service]
ExecStart=/usr/loca/bin/flapperd

[Install]
WantedBy=multi-user.target
